# roles/ansible-letsencrypt/templates/certbot.timer
[Unit]
Description=Twice daily renewal of Let's Encrypt's certificates

[Timer]
OnCalendar=0/12:00:00
RandomizedDelaySec={{RandomizedDelaySec}}
Persistent=true

[Install]
WantedBy=timers.target